# infra/systemd/maistv-worker-media.service
[Unit]
Description=MaisTV BullMQ media worker
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/SaasTV
Environment=NODE_ENV=production
Environment=MAISTV_RUNTIME_ROLE=worker-media
EnvironmentFile=-/root/SaasTV/.env
EnvironmentFile=/etc/maistv/chat-stack.env
ExecStart=/usr/bin/node server/workers/media.worker.mjs
Restart=on-failure
RestartSec=5
TimeoutStopSec=90
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal
SyslogIdentifier=maistv-worker-media

[Install]
WantedBy=multi-user.target
