# infra/systemd/maistv-next-sse.service
[Unit]
Description=MaisTV Next SSE realtime gateway (homologacao)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/MaisTV
Environment=NODE_ENV=production
Environment=SSE_HOST=127.0.0.1
Environment=SSE_PORT=5356
EnvironmentFile=/etc/maistv-next/maistv-next.env
ExecStart=/usr/bin/node server/realtime/sse-server.mjs
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal
SyslogIdentifier=maistv-next-sse

[Install]
WantedBy=multi-user.target
