# infra/systemd/maistv-next-api.service
[Unit]
Description=MaisTV Next local API (homologacao)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/MaisTV
Environment=NODE_ENV=production
Environment=PORT=5353
Environment=MAISTV_RUNTIME_ROLE=chat-api
EnvironmentFile=/etc/maistv-next/maistv-next.env
ExecStart=/usr/bin/node server/local-api.mjs
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal
SyslogIdentifier=maistv-next-api

[Install]
WantedBy=multi-user.target
