# Checkout-Renovacao/deploy/systemd/maistv-worker.service
[Unit]
Description=MaisTV Background Worker
After=network.target maistv-api.service maistv-whatsapp.service

[Service]
Type=simple
WorkingDirectory=/root/SaasTV
EnvironmentFile=/root/SaasTV/.env
Environment=NODE_ENV=production
Environment=MAISTV_RUNTIME_ROLE=worker
Environment=WHATSAPP_HTTP_ENABLED=false
Environment=WHATSAPP_SCHEDULERS_ENABLED=true
ExecStart=/usr/bin/node /root/SaasTV/server/maistv-worker.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
